# aports/main/postmarketos-base/50-firmware.rules
# do not edit this file, it will be overwritten on update
SUBSYSTEM=="firmware", ACTION=="add", RUN+="/usr/bin/firmware.sh"
